# etc/scripts/crontab
59 23 * * SUN sh /etc/scripts/preload_data.sh >> /var/log/cron_complexity.log 2>&1
* * * * * echo "Testing cronjob log: [$(date)]" >> /var/log/test_log.log 2>&1
