# etc/scripts/crontab
59 23 * * SUN sh /salic_ml_web/preload_data.sh >> /var/log/cron_complexity.log 2>&1
